# docker/services/messenger-worker-1/run
#!/usr/bin/with-contenv sh

/opt/run-as-php.sh "/www/bin/console messenger:consume --limit=10 resolve"
