# docker/services/messenger-worker-1/run
#!/bin/sh
with-contenv
exec 2>&1
su php \
    -s "/bin/ash" \
    -c "/www/bin/console messenger:consume --limit=10 resolve"
