# docker/services/doctrine-migrations/up
#!/usr/bin/with-contenv sh

foreground { echo "Doctrine: running database migrations ..." }

/opt/run-as-php.sh "/www/bin/console doctrine:migrations:migrate -n -v"
