# docker/services/doctrine-migrations/up
#!/bin/sh
with-contenv
foreground { echo "Doctrine: running database migrations ..." }

/opt/run-as-php.sh "/www/bin/console doctrine:migrations:migrate -n -v"
