# docker/services/doctrine-migrations/up
foreground { echo "Doctrine: running database migrations ..." }
su php -s "/bin/ash" \
    -c "/www/bin/console doctrine:migrations:migrate -n -v"
