# docker/services/messenger-worker-2/run
#!/bin/sh
with-contenv

/opt/run-as-php.sh "/www/bin/console messenger:consume --limit=10 scheduler_default async"
