# docker/services/messenger-worker-2/run
#!/usr/bin/with-contenv sh
exec 2>&1
su php \
    -s "/bin/ash" \
    -c "/www/bin/console messenger:consume --limit=10 scheduler_default"
